# task-03-server-side-tools/postgres-backup-service/pg-backup.cron
# Run backup script daily at 2:00 AM
0 2 * * * root /usr/local/bin/pg-backup.sh >> /var/log/pg-backup.log 2>&1

